# ansible/roles/debops.opennebula_core/templates/etc/systemd/system/opennebula.service.d/after-mysql.conf.j2
# {{ ansible_managed }}

[Unit]
After=mysql.service
